# config/ansible/roles/cowrie/templates/cowrie.container.j2
# /etc/containers/systemd/users/{{ system.users.sensor.uid }}/cowrie.container
# Managed by Ansible

[Unit]
Description=Cowrie Honeypot Sensor
After=network-online.target

[Container]
Image=docker.io/cowrie/cowrie:latest
ContainerName=sensor-cowrie
Network=pasta

PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:2222:2222
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:2222:2222
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:2223:2223
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:2223:2223

# Hardening
DropCapability=ALL
NoNewPrivileges=true
ReadOnly=true
VolatileTmp=true

# Health check (Podman 5.x Quadlet key names)
HealthCmd=CMD-SHELL echo | nc -w 3 localhost 2222 | grep -q SSH || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

# Resource limits
PodmanArgs=--cpus=0.5 --memory=256m

# Volumes: U for subUID ownership mapping in rootless container
Volume=/etc/lantana/sensor/cowrie/cowrie.cfg:/cowrie/cowrie-git/etc/cowrie.cfg:ro
Volume=/var/log/lantana/sensor/cowrie:/cowrie/cowrie-git/var/log/cowrie:rw,U
Volume=/var/lib/lantana/sensor/cowrie:/cowrie/cowrie-git/var/lib/cowrie:rw,U

[Service]
Restart=always

[Install]
WantedBy=default.target
