# config/ansible/roles/dionaea/templates/dionaea.container.j2
# /etc/containers/systemd/users/{{ system.users.sensor.uid }}/dionaea.container
# Managed by Ansible.
#
# PublishPort lines are generated from dionaea_service_catalog. All
# services are always published — disabling a specific protocol is a
# manual on-server operation. See docs/honeypots.md.

[Unit]
Description=Dionaea Malware Honeypot Sensor
After=network-online.target

[Container]
# :latest is frozen at the Nov-2020 0.11.0 release; :nightly is the
# actively-rebuilt tag (refreshed daily by the upstream CI). Both ship
# the same source codebase (dionaea 0.11.0) but :nightly carries fresh
# OS-package security updates.
Image=docker.io/dinotools/dionaea:nightly
ContainerName=sensor-dionaea
Network=pasta

# Force the entrypoint's init_etc / init_lib / init_log functions to
# run on every container start. Without these, the entrypoint's
# `test ! -d /opt/dionaea/etc/dionaea` check would see our bind-mounted
# subdirectories as existing → skip the seeding → dionaea launches with
# no dionaea.cfg → exit 133 instantly.
Environment=DIONAEA_FORCE_INIT_CONF=1
Environment=DIONAEA_FORCE_INIT_DATA=1

{% for svc, s in dionaea_service_catalog.items() -%}
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:{{ s.pub_port }}:{{ s.int_port }}
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:{{ s.pub_port }}:{{ s.int_port }}
{% if s.transport == 'tcp_udp' -%}
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:{{ s.pub_port }}:{{ s.int_port }}/udp
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:{{ s.pub_port }}:{{ s.int_port }}/udp
{% endif -%}
{% endfor %}

# Hardening
# Capability strategy: drop everything by default, then add back only
# what the bundled image's startup actually needs.
#   NET_BIND_SERVICE — dionaea binds privileged ports (21, 80, 135, 445)
#                     inside the container.
#   SETUID + SETGID  — the entrypoint launches dionaea with `-u dionaea
#                     -g dionaea`; without these the privilege-drop
#                     fails and dionaea exits silently with status 133.
#   CHOWN + FOWNER   — the entrypoint's `cp -a` (template/lib → var/lib)
#                     preserves ownership; without these we get EPERM
#                     noise on every boot. Cosmetic but worth removing.
# ReadOnly=true is incompatible with this image's entrypoint, which has
# to write the seeded dionaea.cfg into /opt/dionaea/etc/dionaea/. The
# bundled dionaea.cfg disables nfq and pcap modules so no CAP_NET_ADMIN
# or CAP_NET_RAW required.
DropCapability=ALL
AddCapability=NET_BIND_SERVICE SETUID SETGID CHOWN FOWNER
NoNewPrivileges=true
VolatileTmp=true

# Health check (Podman 5.x Quadlet key names)
# Process-level: resilient to dionaea_services toggles (port-level checks
# would break if SMB or whichever port we probed got disabled).
HealthCmd=CMD-SHELL pgrep -f dionaea >/dev/null || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

# Resource limits (higher than Cowrie: handles more protocols + stores malware)
PodmanArgs=--cpus=1.0 --memory=512m

# Volumes
#   services-enabled/ and ihandlers-enabled/ are mounted :ro over the
#   bundled directories inside the image. The bundled dionaea.cfg globs
#   *.yaml from these paths and only the files we render here will be
#   loaded — overriding the image's default 16 services (including
#   blackhole on port 23, which would conflict with Cowrie) and 6
#   ihandlers (none of which is log_json).
#   log + lib bind-mounts use U for subUID ownership mapping in rootless
#   container. The entrypoint's init_lib + init_log copy seed dirs
#   (binaries, bistreams, ftp, http, ...) into these mounts on first
#   boot via cp -van (no-clobber).
Volume=/etc/lantana/sensor/dionaea/services-enabled:/opt/dionaea/etc/dionaea/services-enabled:ro
Volume=/etc/lantana/sensor/dionaea/ihandlers-enabled:/opt/dionaea/etc/dionaea/ihandlers-enabled:ro
Volume=/var/log/lantana/sensor/dionaea:/opt/dionaea/var/log/dionaea:rw,U
Volume=/var/lib/lantana/sensor/dionaea:/opt/dionaea/var/lib/dionaea:rw,U

[Service]
Restart=always

[Install]
WantedBy=default.target
