# config/ansible/roles/dionaea/templates/dionaea.container.j2
# /etc/containers/systemd/users/{{ system.users.sensor.uid }}/dionaea.container
# Managed by Ansible

[Unit]
Description=Dionaea Malware Honeypot Sensor
After=network-online.target

[Container]
Image=docker.io/dinotools/dionaea:latest
ContainerName=sensor-dionaea
Network=pasta

PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8021:21
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8021:21
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8080:80
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8080:80
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8135:135
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8135:135
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8445:445
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8445:445
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8433:1433
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8433:1433
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8306:3306
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8306:3306
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8060:5060
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8060:5060
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:8060:5060/udp
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:8060:5060/udp

# Hardening
DropCapability=ALL
NoNewPrivileges=true
ReadOnly=true
VolatileTmp=true

# Health check (Podman 5.x Quadlet key names)
# Process-level: resilient to dionaea_services toggles (port-level checks
# would break if SMB or whichever port we probed got disabled).
HealthCmd=CMD-SHELL pgrep -f dionaea >/dev/null || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

# Resource limits (higher than Cowrie: handles more protocols + stores malware)
PodmanArgs=--cpus=1.0 --memory=512m

# Volumes: U for subUID ownership mapping in rootless container
Volume=/etc/lantana/sensor/dionaea/dionaea.yaml:/opt/dionaea/etc/dionaea/dionaea.yaml:ro
Volume=/var/log/lantana/sensor/dionaea:/opt/dionaea/var/log/dionaea:rw,U
Volume=/var/lib/lantana/sensor/dionaea:/opt/dionaea/var/lib/dionaea:rw,U

[Service]
Restart=always

[Install]
WantedBy=default.target
