# config/ansible/roles/dionaea/templates/dionaea.container.j2
# /etc/containers/systemd/users/{{ system.users.sensor.uid }}/dionaea.container
# Managed by Ansible.
#
# PublishPort lines are generated from dionaea_services × dionaea_service_catalog.
# Disable a service by removing it from dionaea_services in inventory —
# no edits needed here.

[Unit]
Description=Dionaea Malware Honeypot Sensor
After=network-online.target

[Container]
Image=docker.io/dinotools/dionaea:latest
ContainerName=sensor-dionaea
Network=pasta

{% for svc in dionaea_services -%}
{%- set s = dionaea_service_catalog[svc] %}
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:{{ s.pub_port }}:{{ s.int_port }}
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:{{ s.pub_port }}:{{ s.int_port }}
{% if s.transport == 'tcp_udp' -%}
PublishPort={{ network.sensors.sensor_01.ipv4.split('/')[0] }}:{{ s.pub_port }}:{{ s.int_port }}/udp
PublishPort=[{{ network.sensors.sensor_01.ipv6.split('/')[0] }}]:{{ s.pub_port }}:{{ s.int_port }}/udp
{% endif -%}
{% endfor %}

# Hardening
DropCapability=ALL
NoNewPrivileges=true
ReadOnly=true
VolatileTmp=true

# Health check (Podman 5.x Quadlet key names)
# Process-level: resilient to dionaea_services toggles (port-level checks
# would break if SMB or whichever port we probed got disabled).
HealthCmd=CMD-SHELL pgrep -f dionaea >/dev/null || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

# Resource limits (higher than Cowrie: handles more protocols + stores malware)
PodmanArgs=--cpus=1.0 --memory=512m

# Volumes
#   services/ and ihandlers/ are mounted :ro over the bundled directories
#   inside the image — dionaea.cfg globs *.yaml from these paths and only
#   the files we render here will be loaded.
#   log + lib bind-mounts use U for subUID ownership mapping in rootless container.
Volume=/etc/lantana/sensor/dionaea/services:/opt/dionaea/etc/dionaea/services:ro
Volume=/etc/lantana/sensor/dionaea/ihandlers:/opt/dionaea/etc/dionaea/ihandlers:ro
Volume=/var/log/lantana/sensor/dionaea:/opt/dionaea/var/log/dionaea:rw,U
Volume=/var/lib/lantana/sensor/dionaea:/opt/dionaea/var/lib/dionaea:rw,U

[Service]
Restart=always

[Install]
WantedBy=default.target
